# saas/programming-model/builders/specs/deb/compss-framework/etc/profile.d/compss.sh
export PATH=$PATH:/opt/COMPSs/Runtime/scripts/user
